# overlays/firmware-extended/99-remote-screen/root/etc/supervisor.d/remote_screen.ini
[program:remote_screen]
user=remote_screen
command=/usr/local/bin/fb-http-server-run.sh
environment=EXTENDED_CFG="/home/lava/printer_data/config/extended/extended.cfg",
    FB_DEVICE="/dev/fb0",
    TOUCH_DEVICE="/dev/input/event0",
    HTML_TEMPLATE="/usr/local/share/fb-http-server/index.html",
    PORT="8092",
    BIND="127.0.0.1"
autostart=true
autorestart=unexpected
startsecs=1
startretries=3
exitcodes=0
stdout_logfile=/var/log/fb-http-server.log
stderr_logfile=/var/log/fb-http-server.log
redirect_stderr=false
